# deploy/systemd/gille-autonomy-tick.service
[Unit]
Description=gille autonomy controller tick (gi#49)

[Service]
Type=oneshot
WorkingDirectory=/home/magnus/home-server-eval
ExecStart=/usr/bin/env npx tsx scripts/autonomy-tick-cli.ts
