# deploy/systemd/gille-autonomy-tick.timer
[Unit]
Description=Daily gille autonomy tick

[Timer]
OnCalendar=*-*-* 05:30:00
Persistent=true

[Install]
WantedBy=timers.target
